# serial-forwarder.service
[Unit]
Description=Serial to TCP Forwarder - Web Service & Background Forwarder
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/Com2tcp
ExecStart=/path/to/venv/bin/python /path/to/Com2tcp/service_runner.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serial-forwarder

[Install]
WantedBy=multi-user.target
